# .config/systemd/user/dotfiles-backup.timer
[Unit]
Description=Backup dotfiles daily

[Timer]
OnCalendar=daily
Persistent=true
Unit=dotfiles-backup.service

[Install]
WantedBy=timers.target
