# .config/systemd/user/dotfiles-backup.service
[Service]
ExecStart=bash %h/dot/backup.sh

[Unit]
Wants=network-online.target
After=network-online.target nss-lookup.target
