# .config/systemd/user/dotfiles-backup.timer
[Unit]
Description=Backup dotfiles daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=dotfiles-backup.service
